# builder/assets/system/pwnagotchi-syswatchdog.service
[Unit]
Description=pwnagotchi system-level watchdog (recovers from failures the Python process can't see or fix itself)
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/pwnagotchi-syswatchdog
