# builder/data/etc/systemd/system/pwnagotchi-battery-curve-log.service
[Unit]
Description=Log PiSugar 3 battery voltage/percentage for a discharge curve test

[Service]
Type=oneshot
ExecStart=/usr/bin/pwnagotchi-battery-curve-log
